# realtime/snowman-realtime.service
[Unit]
Description=Snowman Realtime Voice Assistant
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=snowman
Group=audio
WorkingDirectory=/home/snowman/voice-assistant-realtime/realtime
EnvironmentFile=/home/snowman/voice-assistant-realtime/realtime/.env
ExecCondition=/home/snowman/voice-assistant-realtime/realtime/within_runtime_window.sh
ExecStart=/home/snowman/voice-assistant-realtime/realtime/start_realtime.sh
Restart=always
RestartSec=5
KillMode=control-group
StandardOutput=append:/home/snowman/voice-assistant-realtime/realtime/realtime.log
StandardError=append:/home/snowman/voice-assistant-realtime/realtime/realtime.log

[Install]
WantedBy=multi-user.target
